# autodialer-stack.target
[Unit]
Description=Autodialer Stack - Bot + Asterisk
After=network.target

[Install]
WantedBy=multi-user.target
