# autodialer.service
[Unit]
Description=Autodialer Pro - Telegram bot + Nonbor buyurtmalar
After=network.target asterisk.service
Wants=asterisk.service
PartOf=autodialer-stack.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/autodialer-pro/src
ExecStart=/usr/bin/python3 autodialer.py
EnvironmentFile=/opt/autodialer-pro/.env
Restart=always
RestartSec=10
StandardOutput=append:/var/log/autodialer/autodialer.log
StandardError=append:/var/log/autodialer/autodialer.log

[Install]
WantedBy=autodialer-stack.target
